# examples/telegram-bot/nexus-health.service
[Unit]
Description=Nexus Health Check Endpoint
After=network.target nexus-bot.service nexus-processor.service
Wants=nexus-bot.service nexus-processor.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/nexus-core/examples/telegram-bot
EnvironmentFile=/opt/nexus/.env
Environment="HEALTH_CHECK_PORT=8080"
Environment="HEALTH_CHECK_HOST=127.0.0.1"
ExecStart=/opt/nexus-core/examples/telegram-bot/venv/bin/python src/health_check.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
