# examples/telegram-bot/nexus-bot.service
[Unit]
Description=Nexus Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/nexus-core/examples/telegram-bot
EnvironmentFile=/opt/nexus/.env
Environment="PYTHONUNBUFFERED=1"
# Include system paths so system binaries (e.g. /usr/bin/gh) are available
Environment="PATH=/opt/nexus-core/examples/telegram-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/nexus-core/examples/telegram-bot/venv/bin/python src/telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
